# config/test_flask.conf
[program:test_flask]
command = /home/python2.7/bin/gunicorn app:app -b 0.0.0.0:8000
directory = /home/test_flask